# packaging/systemd/xhod.service
[Unit]
Description=Cross Host Ops daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/xhod --config /etc/xho/config.toml --origin external
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
